# agh-sync.service
[Unit]
Description=AdGuardHome Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agh-sync
ExecStart=/usr/local/bin/agh-sync run --watch --config /etc/agh-sync/config.yaml
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadOnlyPaths=/etc/agh-sync
ReadWritePaths=/var/lib/agh-sync
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
